# deploy/systemd/euphoric-tickets-restart.timer
[Unit]
Description=Weekly restart of Euphoric Tickets (Tuesday 4 AM)

[Timer]
OnCalendar=Tue *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
